# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:php-server]
command=php -S 0.0.0.0:80
autostart=true
autorestart=true

[program:worker]
command=php worker.php
autostart=true
autorestart=true